# mini-bot/deploy/minibot.timer
[Unit]
Description=Run MiniBot cycle every 4 hours

[Timer]
OnBootSec=5min
OnUnitActiveSec=4h
Persistent=true

[Install]
WantedBy=timers.target
